# launchd/com.horse3.scraping.historique.matin.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.horse3.scraping.historique.matin</string>
    
    <key>ProgramArguments</key>
    <array>
        <string>/Users/gicquelsacha/horse3/cron_scraping_historique.sh</string>
    </array>
    
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>4</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>
    
    <key>StandardOutPath</key>
    <string>/Users/gicquelsacha/horse3/logs/launchd_historique.log</string>
    
    <key>StandardErrorPath</key>
    <string>/Users/gicquelsacha/horse3/logs/launchd_historique_error.log</string>
    
    <key>RunAtLoad</key>
    <false/>
    
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/opt/homebrew/bin</string>
    </dict>
</dict>
</plist>
