# launchd/com.horse.lessons-sync.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.horse.lessons-sync</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/gicquelsacha/horse3/.venv/bin/python</string>
        <string>/Users/gicquelsacha/horse3/web/backend/daily_lessons_sync.py</string>
    </array>

    <!-- Exécution quotidienne à 23h30 -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>23</integer>
        <key>Minute</key>
        <integer>30</integer>
    </dict>

    <!-- Logs -->
    <key>StandardOutPath</key>
    <string>/Users/gicquelsacha/horse3/logs/lessons_sync.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/gicquelsacha/horse3/logs/lessons_sync_error.log</string>

    <!-- Configuration -->
    <key>WorkingDirectory</key>
    <string>/Users/gicquelsacha/horse3</string>

    <!-- Environnement -->
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/Users/gicquelsacha/horse3/.venv/bin</string>
    </dict>

    <!-- Comportement -->
    <key>RunAtLoad</key>
    <false/>

    <key>KeepAlive</key>
    <false/>
</dict>
</plist>
